# docker/etc/supervisor.d/authenticationEngine.ini
[program:authenticationEngine]
process_name:authenticationEngine
autostart=true
autorestart=unexpected
user=root
directory=/home/code/
command=npm start
stderr_logfile = /var/log/supervisord/authenticationEngine-stderr.log
stdout_logfile = /var/log/supervisord/authenticationEngine-stdout.log
